# ansible/installs/files/temperature.service
[Unit]
Description=Temperature

[Service]
User=temperature
Group=temperature
ExecStart=/usr/local/bin/demo

[Install]
WantedBy=multi-user.target
